# contrib/openrc/miraikumiko
#!/sbin/openrc-run

description="miraikumiko backend"

name=$RC_SVCNAME

command="/var/www/$name.space/venv/bin/python /var/www/$name.space/run.py"
command_args="-c /etc/$name.conf -l /var/log/$name.log"
pidfile="/run/$name.pid"

depend() {
    need net
}

start() {
	ebegin "Starting $name backend"
    cd /var/www/$name.space
	start-stop-daemon -bm -S -u root:root -p $pidfile -x $command -- $command_args
	eend $?
}

stop() {
    ebegin "Stopping $name backend"
    start-stop-daemon -K -p $pidfile
    eend $?
}

restart() {
    stop
    sleep 1
    start
}
